# services/load-user-modules/up
/etc/s6-rc/scripts/load-user-modules/up
