# services/set-term-encoding/up
/etc/s6-rc/scripts/set-term-encoding/up
